# ubupdate/etc/systemd/system/ubupdate.service
[Unit]
Description=Ubuntu Auto-Updater (Java)
After=network.target

[Service]
Type=oneshot
ExecStart=java /usr/local/bin/ubupdate.jar